# ultron/tools/log-analysis-api.service
[Unit]
Description=Agent Log Analysis API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/tools
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/tools/log_analysis_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target