# ultron/cost_optimizer/cost-optimizer.service
[Unit]
Description=Ultron Cost Optimizer API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/cost_optimizer
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/cost_optimizer/cost_optimizer.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target