# ultron-workflow/monitoring/alert-history-dashboard.service
[Unit]
Description=Alert History Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron-workflow/monitoring
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron-workflow/monitoring/alert_history_dashboard.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target