# ultron-workflow/agents/queue-api-server.service
[Unit]
Description=Agent Task Queue API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron-workflow/agents
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron-workflow/agents/queue_api_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target