# ultron/agents/agent-api-gateway.service
[Unit]
Description=Agent API Gateway - Unified Entry Point
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/agents
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/agents/agent_api_gateway.py --serve 0.0.0.0 18290
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target