# ultron/agents/agent-dispatcher-api.service
[Unit]
Description=Agent Dispatcher API - Request Load Balancing
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/agents
ExecStart=/usr/bin/python3 agent_dispatcher_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target