# ultron/agents/autoscaler.service
[Unit]
Description=Agent AutoScaler Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/agents
ExecStart=/usr/bin/python3 autoscaler_api.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target