# ultron/agents/collaboration-analytics.service
[Unit]
Description=Collaboration Analytics Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/agents
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/agents/collaboration_analytics.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target