# ultron/agents/result-aggregator.service
[Unit]
Description=Agent Result Aggregator Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/agents
ExecStart=/usr/bin/python3 result_aggregator.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target