# ultron/agent-collaboration-hub.service
[Unit]
Description=Agent Collaboration Hub API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/agent_collaboration_hub.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target