# ultron/alert-intelligence-api.service
[Unit]
Description=Alert Intelligence API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/alert_intelligence_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target