# ultron/alert-integration-api.service
[Unit]
Description=Alert Integration API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/tools/alert_integration_api.py
Restart=always
RestartSec=5

Environment=PYTHONPATH=/root/.openclaw/workspace/ultron/tools

[Install]
WantedBy=multi-user.target