# ultron/circuit-breaker.service
[Unit]
Description=Agent Circuit Breaker and Rate Limiter Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/agent-circuit-breaker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target