# ultron/decision-automation.service
[Unit]
Description=Decision Automation Orchestrator
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron/decision_engine
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/decision_engine/automation.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target