# ultron/ops-metrics-api.service
[Unit]
Description=OpenClaw Ops Metrics API (Enhanced)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/ops-metrics-api.py --port 18220
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target