# ultron/ops-metrics-api.service
[Unit]
Description=奥创运维指标API服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 ops-metrics-api.py --port 18230
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target