# ultron/service-governance-api.service
[Unit]
Description=Ultron Service Governance API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/ultron
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/ultron/service_governance_api.py
Environment="PORT=18250"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target